# pocketchip/config/includes.chroot_before_packages/usr/lib/systemd/system/ubihealthd.service
# DUPLICATE: identical copy in every x-chip-os flavor (headless/gui/pocketchip) -- keep these in sync (no common/ dir; see README).
[Unit]
Description=NAND health watchdog service
ConditionFileIsExecutable=/usr/sbin/ubihealthd
After=systemd-journald.socket

[Service]
Type=simple
ExecStart=/usr/sbin/ubihealthd -d /dev/ubi0 -f /var/cache/ubihealthd.v1.stats
PIDFile=/run/ubihealthd.pid
StandardOutput=journal
Restart=on-failure

[Install]
WantedBy=multi-user.target
